# system/daphne/daphne.socket
[Unit]
Description=daphne daemon

[Socket]
ListenStream=/run/daphne.sock

[Install]
WantedBy=sockets.target
